# overlay/etc/periodic/hourly/cleanup
#!/command/with-contenv sh

set -euo pipefail
export BF_E=`basename ${0}`


#======================================================================================================================
# Use executable to cleanup old uploads.
#======================================================================================================================

send-cleanup
